# dapnet-gateway/files/etc/uci-defaults/99_dapnetgateway-init
#!/bin/sh

ENABLE=$(uci -q get mmdvm.POCSAG.Enable)
if [ "$ENABLE" != "1" ]; then
    /etc/init.d/dapnetgateway stop
    /etc/init.d/dapnetgateway disable
fi

exit 0
