# ysf-clients/files/etc/uci-defaults/99_ysfgateway-init
#!/bin/sh

[ ! -f /etc/config/mmdvm ] && touch /etc/config/mmdvm
[ -z `uci -q get mmdvm.System_Fusion.Enable` ] && {
	uci set mmdvm.System_Fusion=mmdvmhost
	uci set mmdvm.System_Fusion.Enable=1
	uci commit mmdvm
}

ENABLE=$(uci -q get mmdvm.NXDN.Enable)
if [ "$ENABLE" != "1" ]; then
    /etc/init.d/nxdngateway stop && /etc/init.d/nxdngateway disable
    /etc/init.d/nxdnparrot stop && /etc/init.d/nxdnparrot disable
fi

exit 0
